# deploy/macbook/com.example.zotero-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <!-- Rename to com.YOURNAME.zotero-sync.plist and replace all YOURUSER /
         TAILSCALE_HOSTNAME placeholders before installing. -->
    <key>Label</key>
    <string>com.YOURNAME.zotero-sync</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/YOURUSER/Library/Scripts/zotero-sync.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
      <key>REMOTE_HOST</key>
      <string>TAILSCALE_HOSTNAME</string>
      <key>REMOTE_USER</key>
      <string>YOURUSER</string>
      <key>REMOTE_DIR</key>
      <string>/home/YOURUSER/docker-binds/academic-mcp/zotero-mirror</string>
      <key>REMOTE_FT</key>
      <string>/home/YOURUSER/docker-binds/academic-mcp/ft-cache</string>
      <key>ZOTERO_SRC</key>
      <string>/Users/YOURUSER/Zotero</string>
      <!-- Minimum seconds between syncs regardless of WatchPaths trigger rate. -->
      <key>DEBOUNCE_SECONDS</key>
      <string>60</string>
    </dict>

    <!-- Watch the WAL/journal companion only, not the whole ~/Zotero dir.
         These files mtime-bump on real DB writes, not on Spotlight churn.
         Watch storage/ for new ft-cache files added by Zotero's indexer. -->
    <key>WatchPaths</key>
    <array>
      <string>/Users/YOURUSER/Zotero/zotero.sqlite-journal</string>
      <string>/Users/YOURUSER/Zotero/storage</string>
    </array>

    <!-- 30-min fallback; covers wake-from-sleep / network reconnect. -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/YOURUSER/Library/Logs/zotero-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOURUSER/Library/Logs/zotero-sync.err</string>
  </dict>
</plist>
